# src/insomniad.service.in
[Unit]
Description=User space autosleep policy daemon

# Don't start autosuspending the system until it is finished booting
After=multi-user.target

[Service]
User=root
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectKernelModules=true
ProtectKernelTunables=true
ReadWritePaths=/sys/power
CapabilityBoundingSet=CAP_DAC_OVERRIDE
Environment="TIMEOUT=10"
ExecStart=@BINDIR@/@PACKAGE@ -t $TIMEOUT $EXTRA_ARGS

[Install]
WantedBy=multi-user.target
