# src/insomniad.service.in
[Unit]
Description=User space autosleep policy daemon

# Don't start autosuspending the system until it is finished booting
After=multi-user.target

[Service]
Environment="TIMEOUT=10"
ExecStart=@BINDIR@/@PACKAGE@ -t $TIMEOUT

[Install]
WantedBy=multi-user.target
